# tor-bcl.conf
[program:tor-bcl]
command=/usr/bin/python /srv/tor-bcl.py
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/tor-bcl.log
stderr_logfile=/var/log/tor-bcl.log
